# threaddht.service
[Unit]
Description=Thread DHT 

[Service]
ExecStart=/home/pi/unisal/atividades/threadDHT/thread_dht 

[Install]
WantedBy=multi-user.target
